# configs/python38/supervisord.conf
[supervisord]
nodaemon=true
loglevel=warn

[program:set-timezone]
command=/bin/sh -c "dpkg-reconfigure -f noninteractive tzdata"
autostart=true
autorestart=false
